# gxde-first-config.service
[Unit]
Description=Config GXDE

[Service]
Type=forking
User=root
ExecStart=/usr/bin/gxde-config

[Install]
WantedBy=multi-user.target
